# ops/systemd/boonetools-api.service
[Unit]
Description=BooneTools API Service
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/boonetools-backend/backend
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--dns-result-order=ipv4first
EnvironmentFile=/opt/boonetools-backend/backend/.env
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=boonetools-api

[Install]
WantedBy=multi-user.target
