# ops/systemd/rapid-swap-listener.service
[Unit]
Description=Rapid Swap WebSocket Listener
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/boonetools-backend
Environment=NODE_ENV=production
EnvironmentFile=/opt/boonetools-backend/backend/.env
ExecStart=/usr/bin/node scripts/rapid-swap-listener.mjs
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rapid-swap-listener

[Install]
WantedBy=multi-user.target
